# deploy/systemd/serverredus-backend.service
[Unit]
Description=Serverredus Telegram Business Backend
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/opt/serverredus
EnvironmentFile=/opt/serverredus/.env
ExecStart=/opt/serverredus/.venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 8000
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target

